# clothes_photo_bot.service
[Unit]
Description=clothes_photo_bot
After=syslog.target
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/clothes_photo_bot/
ExecStart=/home/clothes_photo_bot/env/bin/python3 bot.py
RestartSec=10
Restart=on-failure

StandardOutput=append:/home/clothes_photo_bot/log_output.log
StandardError=append:/home/clothes_photo_bot/log_error.log

[Install]
WantedBy=multi-user.target
